# .crontab
MAILTO=""
0 4 * * * cd ~/work && ctags -R .
